# ansible/roles/mongodb/templates/resume-mongo-backup.timer.j2
# {{ ansible_managed }}
# Nightly schedule for resume-mongo-backup.service (Resume Builder P1).
[Unit]
Description=Nightly resume-builder Mongo backup

[Timer]
OnCalendar={{ mongodb_backup_oncalendar }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
