# ansible/roles/vault-snapshot/templates/vault-snapshot.timer.j2
# {{ ansible_managed }}
# Daily schedule for vault-snapshot.service (PET-109).
[Unit]
Description=Daily Vault raft snapshot (PET-109)

[Timer]
OnCalendar={{ vault_snapshot_oncalendar }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
