# ansible/roles/agent-loop/templates/automerge-poll.timer.j2
# {{ ansible_managed }}
# Periodic schedule for automerge-poll.service (PET-185). Approved Bucket-A PRs should land
# within a tick of going green, so this runs often. RandomizedDelaySec keeps it from
# coinciding with reviewer-stamp-poll's read-modify-write on the same lock-free bucket
# (automerge only APPENDS an event; the stamp rewrite is the sensitive writer).
[Unit]
Description=Periodic auto-merge of gated Bucket-A worker PRs (PET-185)

[Timer]
OnCalendar={{ agent_loop_automerge_oncalendar }}
Persistent=true
RandomizedDelaySec=180

[Install]
WantedBy=timers.target
