# ansible/roles/agent-loop/templates/engine-loop.service.j2
# {{ ansible_managed }}
# One scheduled tick of the engine (PET-184), triggered by engine-loop.timer. Runs AS the loop
# user so it uses the loop's `claude` login (~/.claude.json, Max plan — no API key), the loop's
# gh config, and the Vault Agent token at ~/.vault-token (engine-run.sh mints the petedio-
# engine[bot] token from kv/services/agent-loop at PUSH time only). engine-loop.sh does at most
# ONE unit of work per tick and exits — its cc-slot flock + pause/off-hours/reviewer guards keep
# the engine boxed and yielding (Pedro > reviewer > engine). See docs/runbooks/engine-loop.md.
[Unit]
Description=One tick of the Co-latro engine loop — Bucket-B kinds (PET-184)
After=network-online.target vault-agent.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ agent_loop_user }}
Group={{ agent_loop_user }}
Environment=HOME={{ agent_loop_home }}
# claude + bun live in the per-user npm prefix, not /usr/bin — include it (and pipx shims).
Environment=PATH={{ agent_loop_engine_path }}
Environment=VAULT_ADDR={{ agent_loop_vault_addr }}
Environment=VAULT_CACERT={{ agent_loop_home }}/.config/vault-agent/vault-ca.crt
# Engine knobs (all overridable via the role defaults). The model is pinned here; effort is the
# thinking-token budget; off-hours/slot-lock/state-home drive the loop's yielding behavior.
Environment=ENGINE_MODEL={{ agent_loop_engine_model }}
Environment=ENGINE_MAX_TURNS={{ agent_loop_engine_max_turns }}
Environment=ENGINE_TIMEOUT_S={{ agent_loop_engine_timeout_s }}
Environment=ENGINE_MAX_THINKING_TOKENS={{ agent_loop_engine_max_thinking_tokens }}
Environment=ENGINE_OFFHOURS={{ agent_loop_engine_offhours }}
Environment=ENGINE_HOME={{ agent_loop_engine_home }}
Environment=ENGINE_STATE_DIR={{ agent_loop_engine_home }}/state
Environment=ENGINE_SLOT_LOCK={{ agent_loop_engine_slot_lock }}
ExecStart=/usr/bin/env bash {{ agent_loop_engine_loop_script }}
NoNewPrivileges=true
