# ansible/roles/agent-loop/templates/fleet-alert.service.j2
# {{ ansible_managed }}
# Push a notification when the fleet needs a human (PET-256). Triggered by
# fleet-alert.timer. Runs AS the loop user for the Vault Agent token (ntfy topic read +
# reviewer[bot] token mint for the open-PR check) and the operator-seeded `mc` alias
# (reads the same MinIO eval logs the fleet page renders). READ-ONLY apart from the
# ntfy POST and its local dedupe state file.
[Unit]
Description=Fleet needs-human push alerts (PET-256)
After=network-online.target vault-agent.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ agent_loop_user }}
Group={{ agent_loop_user }}
Environment=HOME={{ agent_loop_home }}
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=VAULT_ADDR={{ agent_loop_vault_addr }}
Environment=VAULT_CACERT={{ agent_loop_home }}/.config/vault-agent/vault-ca.crt
Environment=FLEET_ALERT_REVIEW_HOURS={{ agent_loop_fleet_alert_review_hours }}
ExecStart=/usr/bin/env bash {{ agent_loop_fleet_alert_script }}
NoNewPrivileges=true
