# ansible/roles/agent-loop/templates/reviewer-stamp-poll.service.j2
# {{ ansible_managed }}
# Auto-stamp pedro_verdict onto closed worker PRs (PET-199). Triggered by
# reviewer-stamp-poll.timer. Runs AS the loop user so it uses the Vault Agent token at
# ~/.vault-token (the script mints the petedio-reviewer[bot] App token via reviewer-mint-token.sh
# for its gh reads) and the `mc` alias the operator seeded for this user (the stamp write goes
# to MinIO).
[Unit]
Description=Auto-stamp pedro_verdict onto closed worker PRs (PET-199)
After=network-online.target vault-agent.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ agent_loop_user }}
Group={{ agent_loop_user }}
Environment=HOME={{ agent_loop_home }}
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# VAULT_ADDR/VAULT_CACERT for the script's GH_TOKEN-from-Vault fallback (the CLI reads the
# renewing token off ~/.vault-token, sunk by vault-agent.service).
Environment=VAULT_ADDR={{ agent_loop_vault_addr }}
Environment=VAULT_CACERT={{ agent_loop_home }}/.config/vault-agent/vault-ca.crt
ExecStart=/usr/bin/env bash {{ agent_loop_stamp_poll_script }}
NoNewPrivileges=true
