# ansible/roles/agent-loop/templates/reviewer-loop.service.j2
# {{ ansible_managed }}
# One scheduled tick of the reviewer auto-launch (PET-247 S1), triggered by reviewer-loop.timer.
# Runs AS the loop user so it uses the loop's `claude` login (~/.claude.json, Max plan), the
# loop's gh config, and the Vault Agent token at ~/.vault-token (reads the READ-ONLY Linear key
# from kv/services/linear for the In-Review gate; reviewer-mint-token.sh mints the petedio-
# reviewer[bot] token at review time only). reviewer-loop.sh reviews at most ONE PR per tick and
# exits — it shares the engine's cc-slot flock and holds REVIEWER_ACTIVE while working, so the
# quota order Pedro > reviewer > engine is enforced. See docs/runbooks/reviewer-loop.md.
[Unit]
Description=One tick of the fleet reviewer auto-launch (PET-247)
After=network-online.target vault-agent.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ agent_loop_user }}
Group={{ agent_loop_user }}
Environment=HOME={{ agent_loop_home }}
# claude lives in the per-user npm prefix, not /usr/bin — include it (and pipx shims).
Environment=PATH={{ agent_loop_engine_path }}
Environment=VAULT_ADDR={{ agent_loop_vault_addr }}
Environment=VAULT_CACERT={{ agent_loop_home }}/.config/vault-agent/vault-ca.crt
# Reviewer knobs (all overridable via the role defaults). The model is pinned here — always the
# Opus tier, asymmetrically STRONGER than the engine's model (review is the stronger gate).
Environment=REVIEWER_MODEL={{ agent_loop_claude_model }}
Environment=REVIEWER_MAX_TURNS={{ agent_loop_reviewer_max_turns }}
Environment=REVIEWER_TIMEOUT_S={{ agent_loop_reviewer_timeout_s }}
Environment=REVIEWER_HOME={{ agent_loop_reviewer_home }}
Environment=REVIEWER_LAUNCH_COOLDOWN_MIN={{ agent_loop_reviewer_launch_cooldown_min }}
Environment=ENGINE_HOME={{ agent_loop_engine_home }}
Environment=ENGINE_SLOT_LOCK={{ agent_loop_engine_slot_lock }}
ExecStart=/usr/bin/env bash {{ agent_loop_reviewer_loop_script }}
NoNewPrivileges=true
