# ansible/roles/agent-loop/templates/worker-loop.service.j2
# {{ ansible_managed }}
# One tick of the worker auto-launch (PET-184 S1), triggered by worker-loop.timer. Runs AS the
# loop user so it uses the loop's gh config, the Vault Agent token at ~/.vault-token (reads the
# READ-ONLY Linear key from kv/services/linear to enumerate candidates, and mints the petedio-
# worker[bot] token at push time), and the local Ollama. worker-loop.sh does at most ONE issue
# per tick (single-instance lock + PAUSED guard). The worker is the FREE local tier, so no
# cc-slot / off-hours. See docs/runbooks/worker-loop.md.
[Unit]
Description=One tick of the Co-latro worker auto-launch (PET-184)
After=network-online.target vault-agent.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ agent_loop_user }}
Group={{ agent_loop_user }}
Environment=HOME={{ agent_loop_home }}
Environment=PATH={{ agent_loop_engine_path }}
Environment=VAULT_ADDR={{ agent_loop_vault_addr }}
Environment=VAULT_CACERT={{ agent_loop_home }}/.config/vault-agent/vault-ca.crt
Environment=WORKER_HOME={{ agent_loop_worker_home }}
Environment=WORKER_LAUNCH_COOLDOWN_MIN={{ agent_loop_worker_launch_cooldown_min }}
ExecStart=/usr/bin/env bash {{ agent_loop_worker_loop_script }}
NoNewPrivileges=true
