# ansible/roles/palworld/templates/palworld.service.j2
# {{ ansible_managed }}
# Palworld 1.0 dedicated server (PET-157). Installed AND (by default) enabled + started
# by the palworld role, so the server auto-starts on container/host boot. Launch day just
# needs the auto-pulled 1.0 update + a restart:
#   sudo -u {{ palworld_user }} {{ palworld_steamcmd_dir }}/steamcmd.sh \
#     +force_install_dir {{ palworld_server_dir }} +login anonymous \
#     +app_update {{ palworld_steam_app_id }} validate +quit
#   systemctl restart {{ palworld_service_name }}
[Unit]
Description=Palworld dedicated server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ palworld_user }}
Group={{ palworld_group }}
WorkingDirectory={{ palworld_server_dir }}
# -useperfthreads/NoAsyncLoadingThread are Pocketpair's recommended perf flags for
# the dedicated server; the port comes from the templated PalWorldSettings.ini.
# -EnableGameDataAPI was here on the earlier PET-278 guess that the /v1/api/game-data
# gate is a launch arg. It is not: the arg was live on the running process while the
# endpoint still 404'd. Neither the ini key nor RCON works either — see the long note
# in the role defaults. Dropped rather than left as a no-op that reads as working.
ExecStart={{ palworld_server_exec }} -useperfthreads -NoAsyncLoadingThread -UseMultithreadForDS
Restart=on-failure
RestartSec=5
# Give the server time to flush saves on stop rather than getting SIGKILLed.
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
