# ansible/playbooks/templates/fleet-evals-mirror.timer.j2
# {{ ansible_managed }}
# Fires the agent-evals mirror ~45s after boot and every ~45s after — within the page's
# ~20s auto-refresh tolerance and the 30-60s freshness target. (PET-187)
[Unit]
Description=Periodic MinIO agent-evals mirror for the fleet site (PET-187)

[Timer]
OnBootSec=45s
OnUnitActiveSec=45s
AccuracySec=5s
Unit=fleet-evals-mirror.service

[Install]
WantedBy=timers.target
