# ansible/playbooks/templates/palworld-panel.service.j2
# {{ ansible_managed }}
[Unit]
Description=Palworld control panel (PET-266) — graceful shutdown/start web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ panel_user }}
Group={{ panel_group }}
WorkingDirectory={{ panel_app_dir }}
EnvironmentFile={{ panel_app_dir }}/backend.env
Environment=PORT={{ panel_port }}
# Loopback ONLY (PET-266). The panel has no auth of its own — Cloudflare Access gates the
# tunnel, and nothing gated the LAN, so :8080 on 0.0.0.0 meant anyone on the play segment
# could power the game server off. The two things that legitimately reach it are a browser
# at the keyboard and the cloudflared connector on THIS host
# (playbooks/configure-palworld-tunnel.yml) — both over 127.0.0.1.
Environment=PANEL_BIND={{ panel_bind }}
# This install's identity. The panel tells console traffic apart from Cloudflare Access web
# users by the presence of an Access header, NOT by address (a co-located connector arrives
# on loopback too) — see the panel repo's backend/src/local.ts. PANEL_ID names the machine
# for the console side so attribution doesn't change if the laptop is renamed.
Environment=PANEL_ID={{ panel_id }}
# email -> character seed. Since the loopback lockdown every phone visitor arrives through
# Cloudflare Access, so the panel already knows their email — this saves them ever being
# asked "who's playing?". Only fills gaps: a person who re-links themselves in the app keeps
# their choice across restarts. Emails must match the Access allow-list in cloudflare-routes.tf.
Environment=PANEL_IDENTITY_MAP={{ panel_identity_map }}
Environment=PUBLIC_DIR={{ panel_app_dir }}/frontend/dist
Environment=PALWORLD_HOST={{ palworld_target_host }}
Environment=PALWORLD_REST_PORT={{ palworld_rest_port }}
Environment=PALWORLD_SSH_HOST={{ palworld_target_host }}
Environment=PALWORLD_SSH_USER=root
Environment=PALWORLD_SSH_KEY={{ panel_app_dir }}/panel_ssh_key
Environment=DEFAULT_COUNTDOWN_MINUTES={{ panel_default_countdown }}
# Session history (PET-266). Pinned to an absolute path under the app dir rather than
# left to default relative to the working directory — the unit's WorkingDirectory could
# change and the history would silently start over in a new file. palpanel owns this
# directory, and SQLite's WAL sidecar files land beside it.
Environment=PANEL_DB={{ panel_db_path }}
ExecStart={{ bun_bin }} {{ panel_app_dir }}/backend/src/index.ts
Restart=on-failure
RestartSec=3
# Light hardening (the service needs outbound to localhost REST + ssh).
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=false

[Install]
WantedBy=multi-user.target
